# upstart/isbn.net.in-job.conf
start on starting isbn.net.in
stop on stopping isbn.net.in

# vim: ft=upstart
